# scripts/auto_expo_start.timer
[Unit]
Description=Start auto expo

[Timer]
Unit=auto_expo.service
#OnCalendar=*-*-* *:35..45:00
#OnCalendar=*-*-* *:55..59:00

# 2019-05-10 to 2019-05-31, 22h~23h every min
OnCalendar=2019-05-10..31 22..23:*:00
# 2019-05-11 to 2019-05-31, 00h~02h every min
OnCalendar=2019-05-11..31 00..02:*:00
OnCalendar=2019-06..08-* 22..23:*:00
OnCalendar=2019-06..08-* 00..02:*:00
#Persistent=true

[Install]
WantedBy=timers.target

# vim: set ft=systemd :
