# scripts/auto_expo_stop.service
[Unit]
Description=Stop auto expo

[Service]
Type=oneshot
ExecStart=systemctl stop auto_expo.service 

# vim: set ft=systemd :
