# scripts/Reboot_period@.timer
[Unit]
Description=Reboot periodically

[Timer]
Unit=systemd-reboot.service
OnBootSec=%i

[Install]
WantedBy=timers.target

# vim: set ft=systemd :
