# scripts/should_i_start.service
[Unit]
Description=Take some frame

[Service]
Type=oneshot
ExecStart=/home/bsst/proj/ASISDKDemo/scripts/should_i_start.py

# vim: set ft=systemd :
